# Dockerfiles/php-worker/supervisord_redis_action.ini
[supervisord]
nodaemon=true
[program:laravel-worker-get_redis_send_data]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/pay.aosheng/artisan queue:work redis --queue=get_redis_insert_mysql --sleep=3 --tries=3 --daemon
stdout_logfile=/var/www/html/pay.aosheng/storage/logs/worker_redis_action.log ; stdout 日志输出位置
logfile_backups=15
autostart=true
autorestart=true
numprocs=3
redirect_stderr=true